# Network_dependent_services.mw
Services like netmount or fetchmail depending on an available network connection could be started/stopped by dhcpcd.<ref name=bug522206>{{Bug|522206}}</ref>  This feature can be implemented with this procedure:

=== Implementation ===

Get the patch for dhcpcd:
{{RootCmd|mkdir -p /etc/portage/patches/net-misc/dhcpcd-6.6.7|wget http://home.arcor.de/charles17/gentoo/99-openrc_dhcpcd_hook.patch -P /etc/portage/patches/net-misc/dhcpcd-6.6.7}}

Get the patch for openrc:
{{RootCmd|mkdir -p /etc/portage/patches/sys-apps/openrc-0.13.11|wget http://home.arcor.de/charles17/gentoo/runscript-background.patch -P /etc/portage/patches/sys-apps/openrc-0.13.11}}

Re-emerge both dhcpcd and openrc:
{{RootCmd|emerge -1avt {{=}}net-misc/dhcpcd-6.6.7 {{=}}sys-apps/openrc-0.13.11}}

Add one line '''start_inactive=true''' to the init script:
{{FileBox|filename=/etc/init.d/dhcpcd|lang=bash|1=
#!/sbin/runscript
# Copyright 2007-2008 Roy Marples <roy@marples.name>
# All rights reserved. Released under the 2-clause BSD license.

command=/sbin/dhcpcd
pidfile=/var/run/dhcpcd.pid
command_args=-q
start_inactive=true
name="DHCP Client Daemon"

depend()
{
        provide net
        need localmount
        use logger network
        after bootmisc modules
        before dns
}
}}

Restart dhcpcd:
{{RootCmd|/etc/init.d/dhcpcd restart}}

Later versions than the presently stable of net-misc/dhcpcd and sys-apps/openrc have not been tested with these patches.

=== Result ===

Services having "need net" in their init.d scripts like fetchmail would then start after dhcpcd is started.
{{FileBox|filename=/etc/init.d/fetchmail|lang=bash|1=
#!/sbin/runscript

piddir=${pid_dir:-/var/run/fetchmail}
pid_file=${piddir}/${RC_SVCNAME}.pid
rcfile=/etc/${RC_SVCNAME}rc

depend() {
        need net
        use mta
}
}}
{{RootCmd|eselect rc start fetchmail
||output=<pre>Starting init script
dhcpcd          * Starting DHCP Client Daemon ...
fetchmail       * WARNING: fetchmail is scheduled to start when dhcpcd has started            [ ok ]
</pre>}}

They will be stopped when dhcpcd turns inactive and will be restarted when dhcpcd is back.
{{Cmd|rc-config show default|output=<pre>
  dhcpcd                    [inactive]
  fetchmail                 [stopped]
</pre>}}

== References ==
{{reflist}}
